# docker/standard/config/mosparo.cron
15 1 * * *   /usr/local/bin/php /mosparo/bin/console mosparo:cleanup-database > /dev/null 2>&1
15 4 * * 5   /usr/local/bin/php /mosparo/bin/console mosparo:geoip2:download-database > /dev/null 2>&1
15 3 * * *   /usr/local/bin/php /mosparo/bin/console mosparo:rule-packages:refresh > /dev/null 2>&1
